# services/pdv-auto-trainer.service
[Unit]
Description=PDV Auto Trainer — treina modelo YOLO automaticamente
After=network-online.target pdv-learning-agent.service

[Service]
Type=oneshot
EnvironmentFile=/etc/pdv-antitheft-agent.env
ExecStart=/usr/bin/python3.8 /opt/pdv-antitheft/pdv_auto_trainer.py
User=root
WorkingDirectory=/opt/pdv-antitheft
StandardOutput=journal
StandardError=journal
TimeoutStartSec=18000

[Install]
WantedBy=multi-user.target
