# services/pdv-antitheft-agent.service
[Unit]
Description=PDV Vision Monitor (Groq API)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/pdv-antitheft-agent.env
ExecStart=/usr/bin/python3.8 /opt/pdv-antitheft/pdv_antitheft_agent.py
Restart=always
RestartSec=10
User=root
WorkingDirectory=/opt/pdv-antitheft

[Install]
WantedBy=multi-user.target
