# services/pdv-visual-alert-worker.service
[Unit]
Description=PDV1 Visual Alert Worker via iMHDX
After=network-online.target pdv-telegram-assistant.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/pdv-telegram-assistant.env
Environment=VISUAL_ALERT_POLL_SECONDS=1
Environment=VISUAL_ALERT_MIN_VALUE=20
Environment=VISUAL_ALERT_SEND_CONFERE=0
Environment=VISUAL_ALERT_START_FRESH=1
Environment=VISUAL_ALERT_DELAY_SECONDS=6
Environment=VISUAL_ALERT_IMAGE_RETRY_SECONDS=8
Environment=VISUAL_ALERT_MAX_IMAGE_ATTEMPTS=6
Environment=VISUAL_PHANTOM_ENABLED=1
Environment=VISUAL_PHANTOM_MOTION_MEAN_MIN=4.5
Environment=VISUAL_PHANTOM_CHANGED_PIXELS_MIN=7.0
Environment=VISUAL_IMHDX_LIVE_BUFFER_ENABLED=1
Environment=VISUAL_IMHDX_LIVE_BUFFER_FPS=2
ExecStart=/usr/bin/python3.8 /opt/pdv-visual-auditor/pdv_visual_alert_worker.py
Restart=always
RestartSec=5
User=root
WorkingDirectory=/opt/pdv-visual-auditor

[Install]
WantedBy=multi-user.target
