# services/pdv-auto-trainer.service
[Unit]
Description=PDV Auto Trainer — treina modelo YOLO automaticamente
After=network-online.target pdv-learning-agent.service

[Service]
Type=oneshot
EnvironmentFile=/etc/pdv-antitheft-agent.env
# nice 19 = menor prioridade de CPU possível (não compete com o caixa)
# ionice -c 3 = só usa disco quando ninguém mais precisa
ExecStart=/usr/bin/nice -n 19 /usr/bin/ionice -c 3 /usr/bin/python3.8 /opt/pdv-antitheft/pdv_auto_trainer.py
User=root
WorkingDirectory=/opt/pdv-antitheft
StandardOutput=journal
StandardError=journal
# Limita a 40% do CPU total — caixa sempre tem prioridade
CPUQuota=40%
# Limita RAM a 3GB — não pressiona o restante do sistema
MemoryLimit=3G
# Se demorar mais de 5h, cancela (proteção)
TimeoutStartSec=18000

[Install]
WantedBy=multi-user.target
