# services/pdv-auto-trainer.timer
[Unit]
Description=PDV Auto Trainer — executa a cada hora

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true
Unit=pdv-auto-trainer.service

[Install]
WantedBy=timers.target
